# image_files/lib/systemd/system/upnp.timer
[Timer]
OnStartupSec=5
OnUnitActiveSec=30min

[Install]
WantedBy=multi-user.target
